# etc/cron.d/neumatic_lease_reaper
# Get the chef status for all neumatic servers
0 6 * * * stsapps /var/www/html/neumatic/module/Neumatic/bin/lease_reaper.php > /var/www/html/neumatic/log/lease_reaper.log 2>&1
