# etc/cron.d/neumatic_get_chef_status
# Get the chef status for all neumatic servers
* * * * * stsapps /var/www/html/neumatic/module/Neumatic/bin/get_chef_status.php > /var/www/html/neumatic/log/get_chef_status.log 2>&1
